# templates/gitlab_sidekiq.erb
#!/bin/sh
#
# This file is managed via Puppet

daemon="bundle<%= @ruby_suffix %> exec sidekiq -d"
daemon_flags="-C <%= @sidekiq_config %> -L <%= @sidekiq_log %> -e production"
daemon_user="<%= @gitlab_user %>"
daemon_timeout="300"

. /etc/rc.d/rc.subr

pexp="ruby<%= @ruby_suffix %>:.*sidekiq.*gitlab.*"
rc_reload=NO

rc_start() {
	install -d -o ${daemon_user} -m <%= @gitlab_rundir_mode %> /var/run/gitlab
	${rcexec} "cd <%= @unicorn_root %> && \
		${daemon} ${daemon_flags}"
}

rc_cmd $1
