# templates/gitlab_workhorse.erb
#!/bin/sh
#
# This file is managed via Puppet

daemon="gitlab-workhorse"
daemon_flags="-listenUmask 0 -listenNetwork unix -listenAddr <%= @workhorse_socket %> -authBackend http://127.0.0.1:<%= @gitlab_port %><%= @gitlab_relative_web_path %> -authSocket <%= @gitlab_socket %> -documentRoot <%= @workhorse_document_root %>"
daemon_user="<%= @gitlab_user %>"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_start() {
	install -d -o ${daemon_user} /var/run/gitlab
	${rcexec} "cd <%= @workhorse_root %> && \
		/usr/bin/env PATH=<%= @workhorse_root %>:$PATH ${daemon} ${daemon_flags} &"
}

rc_cmd $1
