# templates/gitlab_unicorn.erb
#!/bin/sh
#
# This file is managed via Puppet

daemon="bundle<%= @ruby_suffix %> exec unicorn_rails<%= @ruby_suffix %> -D"
daemon_flags="-c <%= @unicorn_root %>/config/unicorn.rb -E production --path <%= @unicorn_root %>"
daemon_user="<%= @gitlab_user %>"
daemon_timeout=300

. /etc/rc.d/rc.subr

pexp="ruby<%= @ruby_suffix %>: unicorn_rails<%= @ruby_suffix %> master -D ${daemon_flags:+${daemon_flags}}+(ruby<%= @ruby_suffix %>)*"
rc_reload=NO

rc_start() {
	install -d -o ${daemon_user} /var/run/gitlab
	${rcexec} "cd <%= @unicorn_root %> && \
		${daemon} ${daemon_flags} &"
}

rc_cmd $1
